# scripts/task_blogmaker.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  blogMaker 매일 1회 사이클.
  이 XML 은 참고용 템플릿이다. 환경 별 차이(사용자 SID, Python 절대경로 등)가 있어
  실제 등록은 setup_task.ps1 사용을 권장한다.

  수동 import 시:
    schtasks /create /xml scripts\task_blogmaker.xml /tn "BlogMaker\Cycle"
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>blogMaker 자동 사이클 — 트렌드 수집 → 클러스터링 → 글 생성 → J-Blog 발행</Description>
    <Author>blogMaker</Author>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-01-01T09:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>true</WakeToRun>
    <ExecutionTimeLimit>PT30M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -File "C:\Users\김은희\Downloads\blogMaker\scripts\run_cycle.ps1"</Arguments>
      <WorkingDirectory>C:\Users\김은희\Downloads\blogMaker</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
